# package/systemd/idp.service
[Unit]
Description=Idp service
After=oxauth.service

[Service]
Type=forking
ExecStart=/opt/dist/scripts/idp start
ExecStop=/opt/dist/scripts/idp stop
User=jetty
Group=jetty
TimeoutSec=300

[Install]
WantedBy=multi-user.target
